# devservices/programs.conf
[program:devserver]
command=sentry devserver
autostart=false

[program:taskworker]
command=sentry run taskworker --autoreload
autostart=false
autorestart=true

[program:taskworker-scheduler]
command=sentry run taskworker-scheduler
autostart=false
autorestart=true
