# devservices/programs.conf
[program:devserver]
command=sentry devserver
autostart=false

[program:taskworker]
command=sentry run taskworker --autoreload
autostart=false
autorestart=true

[program:taskworker-scheduler]
command=sentry run taskworker-scheduler
autostart=false
autorestart=true

[program:ingest-events]
command=sentry run consumer ingest-events --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-attachments]
command=sentry run consumer ingest-attachments --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-transactions]
command=sentry run consumer ingest-transactions --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-monitors]
command=sentry run consumer ingest-monitors --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:monitors-clock-tick]
command=sentry run consumer monitors-clock-tick --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:monitors-clock-tasks]
command=sentry run consumer monitors-clock-tasks --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:monitors-incident-occurrences]
command=sentry run consumer monitors-incident-occurrences --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-profiles]
command=sentry run consumer ingest-profiles --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-occurrences]
command=sentry run consumer ingest-occurrences --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:process-spans]
command=sentry run consumer process-spans --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false

[program:process-segments]
command=sentry run consumer process-segments --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:uptime-results]
command=sentry run consumer uptime-results --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-metrics]
command=sentry run consumer ingest-metrics --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-generic-metrics]
command=sentry run consumer ingest-generic-metrics --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:billing-metrics-consumer]
command=sentry run consumer billing-metrics-consumer --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:post-process-forwarder-errors]
command=sentry run consumer post-process-forwarder-errors --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:post-process-forwarder-transactions]
command=sentry run consumer post-process-forwarder-transactions --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:post-process-forwarder-issue-platform]
command=sentry run consumer post-process-forwarder-issue-platform --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:ingest-feedback-events]
command=sentry run consumer ingest-feedback-events --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:eap-spans-subscription-results]
command=sentry run consumer eap-spans-subscription-results --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:subscription-results-eap-items]
command=sentry run consumer subscription-results-eap-items --consumer-group=sentry-consumer --auto-offset-reset=latest --no-strict-offset-reset
autostart=false
autorestart=true

[program:worker]
command=sentry run worker -c 1 --autoreload
autostart=false
autorestart=true
