# devservices/programs.conf
[program:devserver]
command=sentry devserver
autostart=false
